# src/ubuntu/etc/supervisor/conf.d/consul.conf
[program:consul]
command=consul agent -config-dir /etc/consul.d -join %(ENV_VM_IPADDR)s
stopsignal=INT
priority=1
