# src/ubuntu/etc/supervisor/conf.d/cron.conf
[program:cron]
command=/usr/sbin/cron -f
priority=1
