# src/ubuntu/etc/supervisor/conf.d/rsyslog.conf
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
redirect_stderr=true
priority=1
stopsignal=TERM
