# src/graphite/etc/supervisor/conf.d/carbon.conf
[program:carbon-cache]
command=/usr/bin/carbon-cache --config=/etc/carbon/carbon.conf --logdir=/var/log/carbon/ --debug start
priority=10
