# src/rabbitmq/etc/supervisor/conf.d/rabbitmq.conf
[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
user=rabbitmq
