# src/statsd/etc/supervisor/conf.d/statsd.conf
[program:statsd]
command=/usr/local/bin/statsd
