# src/influxdb/etc/supervisor/conf.d/influxdb.conf
[program:infuxdb]
command=/opt/influxdb/current/influxdb -config=/opt/influxdb/shared/config.toml
priority=10
