# src/grafana/etc/supervisor/conf.d/config.conf
[program:config]
command=/usr/local/bin/config.sh
priority=20
autorestart=false
startretries=0
