# src/couchbase/etc/supervisor/conf.d/couchbase.conf
[program:couchbase]
command=/opt/couchbase/bin/couchbase-server -- -noinput
user=couchbase
priority=10
