# src/couchbase/etc/supervisor/conf.d/initialize.conf
[program:initialize]
command=/usr/local/bin/initialize.sh
priority=20
autorestart=false
startretries=0
