# src/memcached/memcached.conf
[program:memcached]
command=/usr/bin/memcached
user=memcache
priority=10
